# deploy/miniscada-kiosk.service
[Unit]
Description=Mini SCADA Local Kiosk Display
After=graphical.target network-online.target miniscada-web.service
Requires=miniscada-web.service
Wants=network-online.target graphical.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=/etc/default/miniscada
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/__APP_USER__/.Xauthority
ExecStartPre=/bin/sleep 8
ExecStart=__APP_DIR__/scripts/start_kiosk.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
