# deploy/miniscada-web.service
[Unit]
Description=Mini SCADA Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=/etc/default/miniscada
ExecStart=__APP_DIR__/.venv/bin/gunicorn -w 2 -b 0.0.0.0:__WEB_PORT__ web.app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
