# deploy/miniscada-web.service
[Unit]
Description=Mini SCADA Web
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/miniscadaModbus
Environment="MINISCADA_SECRET_KEY=change-me"
ExecStart=/opt/miniscadaModbus/.venv/bin/gunicorn -w 2 -b 0.0.0.0:8000 web.app:app
Restart=always

[Install]
WantedBy=multi-user.target
