# deploy/miniscada-simulator.service
[Unit]
Description=Mini SCADA PM8000 Simulator Fleet
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=/etc/default/miniscada
ExecStart=__APP_DIR__/.venv/bin/python -m simulator.pm8000_fleet
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
