# deploy/miniscada-kiosk.service
[Unit]
Description=Mini SCADA Local Kiosk Display
After=network-online.target miniscada-web.service
Requires=miniscada-web.service
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=/etc/default/miniscada
Environment=DISPLAY=:0
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal
ExecStart=/usr/bin/xinit __APP_DIR__/scripts/start_kiosk.sh -- :0 vt1 -keeptty -nocursor
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
